# dovecot-jwks.service
[Unit]
Description=Dovecot Dict Proxy (JWKS key lookup)
Documentation=https://doc.dovecot.org/developer_manual/design/dict_protocol/
# Dovecot depends on this service — start us first, stop us last
Before=dovecot.service
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/dovecot-jwks

# systemd creates these directories and sets the corresponding env vars:
#   CONFIGURATION_DIRECTORY=/etc/dovecot-jwks  (config file location)
#   STATE_DIRECTORY=/var/lib/dovecot-jwks      (output file location)
#   RUNTIME_DIRECTORY=/run/dovecot-jwks        (Unix socket location)
ConfigurationDirectory=dovecot-jwks
StateDirectory=dovecot-jwks
RuntimeDirectory=dovecot-jwks

Restart=on-failure
RestartSec=5s

# Run as the dovecot user so the Unix socket is accessible to Dovecot
User=dovecot
Group=dovecot

# Harden the service — ProtectSystem=strict is safe because all writable
# paths are managed by the *Directory directives above.
PrivateTmp=true
ProtectSystem=strict
NoNewPrivileges=true

[Install]
# Pulling this unit into dovecot.service's dependency graph causes systemd to
# start it automatically whenever Dovecot is started or enabled.
RequiredBy=dovecot.service